# net-mgmt/victoriametrics/files/vmagent.in
#!/bin/sh
# $FreeBSD$
#
# PROVIDE: vmagent
# REQUIRE: LOGIN NETWORKING
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable vmagent:
#
# vmagent_enable (bool):          Set to "NO" by default.
#                                 Set it to "YES" to enable vmagent.
# vmagent_flags (str):            Extra flags to pass to vmagent
#                                 Default: empty

. /etc/rc.subr

name=vmagent
rcvar=vmagent_enable
desc="vmagent daemon"

load_rc_config ${name}

: ${vmagent_enable:=NO}

pidfile="/var/run/${name}.pid"
command="%%PREFIX%%/bin/${name}"
start_cmd="/usr/sbin/daemon -f -p ${pidfile} ${command_interpreter} ${command} ${vmagent_flags}"

run_rc_command "$1"
