# net-mgmt/victoria-metrics/files/vmagent.in
#!/bin/sh
# $FreeBSD$
#
# PROVIDE: vmagent
# REQUIRE: LOGIN NETWORKING
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable vmagent:
#
# vmagent_enable (bool):          Set to "NO" by default.
#                                 Set it to "YES" to enable vmagent.
# vmagent_args (str):             Extra flags to pass to vmagent.
#                                 Default: -enableTCP6
#                                 Example: -remoteWrite.url http://victoriametrics.test:8428/api/v1/write

. /etc/rc.subr

name=vmagent
rcvar=vmagent_enable
desc="vmagent daemon"

load_rc_config ${name}

: ${vmagent_enable:=NO}
: ${vmagent_args:="-enableTCP6"}

vmagent_user=victoria-metrics
vmagent_group=victoria-metrics
vmagent_chdir=/var/db/victoria-metrics/

pidfile="/var/run/${name}.pid"
command="%%PREFIX%%/bin/${name}"
start_cmd="/usr/sbin/daemon -f -p ${pidfile} -u ${vmagent_user} ${command} ${vmagent_args}"

run_rc_command "$1"
